# deploy/portfolio.service
[Unit]
Description=SSH Portfolio TUI
After=network.target
Wants=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/portfolio
ExecStart=/opt/portfolio/portfolio
Environment=PORTFOLIO_ADDR=0.0.0.0:23234
Environment=PORTFOLIO_HOST_KEY=/opt/portfolio/host_key
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
